# conf/scraper_monitor_sup.conf
[program:scraper_monitor]
command=/usr/local/bin/uwsgi --ini /var/www/scraper_monitor/conf/scraper_monitor.ini